# autotag.timer
[Unit]
Description=Run autotag every day at 03:00

[Timer]
OnCalendar=*-*-* 3:00:00

[Install]
WantedBy=timers.target